# project_examples/django_example/templates/golive/cron/db_host.crontab
{% extends "golive/cron/base_db_host.crontab" %}
{% block crontab %}
7 11 * * * {{ USER }} sudo su - postgres -c "/usr/bin/vacuumdb --verbose --all --analyze" >> {{ LOGDIR }}/postgres.log
{% endblock %}
